# crontab
*/10 * * * * /usr/bin/osascript YOUR_SRC_DIR/iTunes-Song-Rating/lib/ranking_progress.scpt
